# django-landysh.service
[Unit]
Description=Django Landysh Application
After=network.target

[Service]
Type=simple
User=root

# Критически важные настройки локали для корректной работы с RAC
# Устанавливаем русскую локаль UTF-8 для всех категорий
Environment="LANG=ru_RU.UTF-8"
Environment="LC_ALL=ru_RU.UTF-8"
Environment="LC_CTYPE=ru_RU.UTF-8"
Environment="LC_MESSAGES=ru_RU.UTF-8"
Environment="LC_NUMERIC=ru_RU.UTF-8"
Environment="LC_TIME=ru_RU.UTF-8"
Environment="LC_COLLATE=ru_RU.UTF-8"
Environment="LC_MONETARY=ru_RU.UTF-8"
Environment="LC_PAPER=ru_RU.UTF-8"
Environment="LC_NAME=ru_RU.UTF-8"
Environment="LC_ADDRESS=ru_RU.UTF-8"
Environment="LC_TELEPHONE=ru_RU.UTF-8"
Environment="LC_MEASUREMENT=ru_RU.UTF-8"
Environment="LC_IDENTIFICATION=ru_RU.UTF-8"

# Настройки для Python - принудительно UTF-8
Environment="PYTHONIOENCODING=utf-8"
Environment="PYTHONUTF8=1"
Environment="PYTHONLEGACYWINDOWSSTDIO=utf-8"

# Настройки для корректной работы с subprocess и RAC
# Убеждаемся, что все команды выполняются с правильной локалью
Environment="LANGUAGE=ru_RU:ru"
Environment="CHARSET=UTF-8"

# Рабочая директория
WorkingDirectory=/data/landysh

# Команда запуска
ExecStart=/data/landysh/venv/bin/python manage.py runserver 0.0.0.0:8000

# Перезапуск при сбое
Restart=on-failure
RestartSec=5

# Лимиты ресурсов (опционально)
# LimitNOFILE=65536

# Безопасность (опционально, если нужно запускать не от root)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target


